# templates/launchagents/com.nexo.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.watchdog
     Runs nexo-watchdog.sh every 30 minutes to monitor system health.
     Checks that the NEXO MCP server, cron jobs, and key infrastructure
     services are alive. Writes results to watchdog-status.json and
     watchdog-report.txt. Sets a .watchdog-alert flag when failures are
     detected so the next Claude session can alert the user immediately.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{NEXO_HOME}}/scripts/nexo-watchdog.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/logs/watchdog-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/logs/watchdog-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NEXO_HOME</key>
        <string>{{NEXO_HOME}}</string>
        <key>NEXO_CODE</key>
        <string>{{NEXO_CODE}}</string>
    </dict>
</dict>
</plist>
