# templates/launchagents/com.nexo.evolution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.evolution
     Static fallback template for the weekly evolution cycle.
     Managed installs now derive a machine-staggered weekly slot from the
     cron manifest to avoid bunching public evolution PRs on Sunday.
     If you install this plist manually, customize Weekday/Hour/Minute
     yourself instead of assuming the managed stagger applies here.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.evolution</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{NEXO_HOME}}/scripts/nexo-evolution-run.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/logs/evolution-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/logs/evolution-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>NEXO_HOME</key>
        <string>{{NEXO_HOME}}</string>
        <key>NEXO_CODE</key>
        <string>{{NEXO_CODE}}</string>
    </dict>
</dict>
</plist>
