# templates/launchagents/com.nexo.immune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.immune
     Runs nexo-immune.py every 30 minutes to perform immune-system checks
     on the NEXO knowledge base. Detects contradictions, stale data, and
     memory drift. Part of NEXO's self-maintenance layer — ensures the
     cognitive database stays internally consistent over time.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.immune</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{NEXO_HOME}}/scripts/nexo-immune.py</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/coordination/immune-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/coordination/immune-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>NEXO_HOME</key>
        <string>{{NEXO_HOME}}</string>
        <key>NEXO_CODE</key>
        <string>{{NEXO_CODE}}</string>
    </dict>
</dict>
</plist>
