# templates/launchagents/com.nexo.postmortem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.postmortem
     Runs nexo-postmortem-consolidator.py every day at 23:30 to produce
     an end-of-day postmortem. Consolidates the day's session diaries,
     decisions, changes, and errors into a single structured summary.
     This summary feeds into the next morning's synthesis and provides
     a permanent record of what happened each day. Does not run at load
     — only fires at the scheduled time.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.postmortem</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{NEXO_HOME}}/scripts/nexo-postmortem-consolidator.py</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/logs/postmortem-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/logs/postmortem-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
