# templates/launchagents/com.nexo.followup-hygiene.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.followup-hygiene
     Runs nexo-followup-hygiene.py every Sunday at 05:00 to clean up
     stale followups and reminders. Archives or closes items that have
     been pending for too long without activity, deduplicates near-
     identical entries, and generates a weekly hygiene report. Keeps the
     operational database from accumulating noise that would clutter
     NEXO's startup menu and priority queue.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.followup-hygiene</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{NEXO_HOME}}/scripts/nexo-followup-hygiene.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/logs/followup-hygiene-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/logs/followup-hygiene-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
