# templates/launchagents/com.nexo.evolution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.evolution
     Runs nexo-evolution-run.py every Sunday at 05:00 to perform a weekly
     evolution cycle. Reviews patterns from the past week, proposes
     improvements to NEXO's own behavior and configuration, and stages
     them as pending proposals for the user to approve. This is NEXO's
     self-improvement mechanism — it cannot apply changes autonomously,
     only propose them.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.evolution</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{NEXO_HOME}}/scripts/nexo-evolution-run.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/logs/evolution-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/logs/evolution-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>NEXO_HOME</key>
        <string>{{NEXO_HOME}}</string>
        <key>NEXO_CODE</key>
        <string>{{NEXO_CODE}}</string>
    </dict>
</dict>
</plist>
