# templates/launchagents/com.nexo.self-audit 2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.nexo.self-audit
     Runs nexo-daily-self-audit.py every day at 07:00 to assess the
     health of NEXO's learning system. Calculates the learning repetition
     rate (how often new learnings duplicate existing ones), identifies
     the top areas with recurring errors, and flags learnings that have
     been ignored 3+ times. Alerts if the repetition rate exceeds 30%.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexo.self-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{NEXO_HOME}}/scripts/nexo-daily-self-audit.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{NEXO_HOME}}/logs/self-audit-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{NEXO_HOME}}/logs/self-audit-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:__NEXO_HOME__/bin</string>
        <key>NEXO_HOME</key>
        <string>{{NEXO_HOME}}</string>
        <key>NEXO_CODE</key>
        <string>{{NEXO_CODE}}</string>
    </dict>
</dict>
</plist>
